# systemd/ghostpi.service
[Unit]
Description=GhostPi Wireless Pentest Tool
Documentation=file:///home/admin/ghostpi/README.md
# network.target: SSH and basic network services are up before we start.
# Deliberately NOT waiting on network-online.target — on a Pi Zero with only
# USB gadget (usb0) that target can take 90 s to time out, delaying the
# e-ink boot prompt and the operator's ability to SSH in.
After=network.target

[Service]
Type=simple
# Root required for raw WiFi socket access (monitor mode, airodump-ng)
User=root
Group=root

WorkingDirectory=/home/admin/ghostpi
ExecStart=/usr/bin/python3 /home/admin/ghostpi/ghostpi/main.py

# Give airmon-ng time to clean up on stop before SIGKILL
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGTERM

# Restart policy: always restart, but wait 5 s to avoid rapid-restart loops
Restart=always
RestartSec=5

# stdout/stderr go to journald; also written to ~/ghostpi/logs/ghostpi.log
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ghostpi

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/admin/ghostpi/ghostpi

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
