# systemd/ghostpi.service
[Unit]
Description=GhostPi Wireless Pentest Tool
Documentation=file:///home/admin/ghostpi/README.md
After=network.target network-online.target
Wants=network-online.target
# Give the USB gadget network interface time to come up
After=systemd-networkd.service

[Service]
Type=simple
# Root required for raw WiFi socket access (monitor mode, airodump-ng)
User=root
Group=root

WorkingDirectory=/home/admin/ghostpi
ExecStart=/usr/bin/python3 /home/admin/ghostpi/ghostpi/main.py

# Give airmon-ng time to clean up on stop before SIGKILL
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGTERM

# Restart policy: always restart, but wait 5 s to avoid rapid-restart loops
Restart=always
RestartSec=5

# stdout/stderr go to journald; also written to ~/ghostpi/logs/ghostpi.log
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ghostpi

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/admin/ghostpi/ghostpi

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
